# s6-service-supervisor/etc/services.d/nginx/run
#!/usr/bin/execlineb -P

# https://github.com/gliderlabs/docker-alpine/issues/185
nginx -g "daemon off;"